# deploy/bore-tunnel.service
[Unit]
Description=Public TCP tunnel for the ascii-royale SSH arena (bore.pub)
After=network-online.target ssh.service
Wants=network-online.target

[Service]
# boxd has no raw-TCP ingress (HTTPS-only edge), so the public SSH port
# rides a bore.pub tunnel: ssh -p 22222 play@bore.pub -> VM port 2222.
# If boxd ever ships TCP port forwarding, delete this unit and publish
# the boxd address instead.
ExecStart=/usr/local/bin/bore local 2222 --to bore.pub --port 22222
Restart=always
RestartSec=5
User=boxd

[Install]
WantedBy=multi-user.target
